# deploy/iphone-store.service
[Unit]
Description=iPhone Store (Next.js)
After=network.target mysql.service
Wants=mysql.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/iphone-store
Environment=NODE_ENV=production
EnvironmentFile=/var/www/iphone-store/.env
ExecStart=/usr/bin/npm run start
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
UMask=0002

[Install]
WantedBy=multi-user.target
